# focal-rt-ros2/rootfs/etc/profile.d/99-source-ros.sh
if [ -f /opt/ros/galactic/setup.bash ]; then
  source /opt/ros/galactic/setup.bash
fi
