# image_builder/data/jammy-rt-iron/rootfs/etc/profile.d/99-source-ros.sh
if [ -f /opt/ros/iron/setup.bash ]; then
  source /opt/ros/iron/setup.bash
fi
