# image_builder/data/noble-rt-jazzy/rootfs/etc/profile.d/99-source-ros.sh
if [ -f /opt/ros/jazzy/setup.bash ]; then
  source /opt/ros/jazzy/setup.bash
fi
