# image_builder/data/jammy-rt-rolling-stanley/rootfs/etc/profile.d/99-source-ros.sh
if [ -f /opt/ros/rolling/setup.bash ]; then
  source /opt/ros/rolling/setup.bash
fi
if [ -f /opt/stanley2/stanley2_control/setup.bash ]; then
  source /opt/stanley2/stanley2_control/setup.bash
fi
if [ -f /opt/stanley_ws/setup.bash ]; then
  source /opt/stanley_ws/setup.bash
fi