# image_builder/data/focal-rt-noetic/rootfs/etc/profile.d/99-source-ros.sh
if [ -f /opt/ros/noetic/setup.bash ]; then
  source /opt/ros/noetic/setup.bash
fi
