# image_builder/data/jammy-rt-rolling/rootfs/etc/profile.d/99-source-ros.sh
if [ -f /opt/ros/rolling/setup.bash ]; then
  source /opt/ros/rolling/setup.bash
fi
